# backup/vaultwarden-prune.timer
[Unit]
Description=Weekly Vaultwarden prune timer

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
