# backup/vaultwarden-backup.timer
[Unit]
Description=Daily Vaultwarden backup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
